# docs/build/tukuyomi.service.example
[Unit]
Description=tukuyomi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tukuyomi
Group=tukuyomi
WorkingDirectory=/opt/tukuyomi
EnvironmentFile=/etc/tukuyomi/tukuyomi.env
ExecStart=/opt/tukuyomi/bin/tukuyomi
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=35
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
